# supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
stdout_logfile=/app/logs/nginx_stdout.log
stderr_logfile=/app/logs/nginx_stderr.log

[program:backend]
command=/app/venv/bin/python /app/backend/main.py
directory=/app/backend
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
stdout_logfile=/app/logs/backend_stdout.log
stderr_logfile=/app/logs/backend_stderr.log
environment=PYTHONPATH="/app"

[program:api]
command=/app/api/target/release/api
directory=/app/api
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
stdout_logfile=/app/logs/api_stdout.log
stderr_logfile=/app/logs/api_stderr.log

[program:web]
command=/app/web/target/release/web
directory=/app/web
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=5
stdout_logfile=/app/logs/web_stdout.log
stderr_logfile=/app/logs/web_stderr.log

